# scripts/production/systemd/user/refinery.slice
# SPDX-FileCopyrightText: 2024 The Refinery Authors <https://refinery.tools/>
#
# SPDX-License-Identifier: EPL-2.0

[Unit]
Description=Groups services related to Refinery

[Slice]
CPUWeight=100
